# .config/systemd/user/dunst.service.d/override.conf
[Unit]
After=startx.service dbus.service
Requires=startx.service dbus.service
