# .local/share/systemd/user/drive.service
[Unit]
Description=Run google drive sync

[Service]
ExecStart=/opt/work/backup/drive/drive_sync.sh
