# .local/share/systemd/user/drive.timer
[Unit]
Description=Sync google drive every hour

[Timer]
OnBootSec=1min
OnUnitActiveSec=1h
Unit=drive.service

[Install]
WantedBy=default.target