# .local/share/systemd/user/grive.service
[Unit]
Description=Run google drive sync

[Service]
Type=forking
ExecStart=/media/data/work/backup/drive/grive_sync.sh
