# .local/share/systemd/user/rclone.timer
[Unit]
Description=Sync google drive every hour

[Timer]
OnBootSec=1min
OnUnitActiveSec=1h
Unit=rclone.service

[Install]
WantedBy=default.target