# .local/share/systemd/user/rclone.service
[Unit]
Description=Run google drive sync

[Service]
ExecStart=/opt/work/backup/drive/rclone_sync.sh
