# .local/share/systemd/user/sleep.service
[Unit]
Description=Run lock screen

[Service]
Type=forking
ExecStart=/home/crandel/.sleep.fish
