# .local/share/systemd/user/sleep.service
[Unit]
Description=Run lock screen

[Service]
Type=Simple
ExecStart=/home/crandel/.sleep.sh
