# .local/share/systemd/user/sleep.service
[Unit]
Description=Run lock screen

[Service]
ExecStart=/home/crandel/.local/bin/.sleep.sh
