# pacmanCache/docker-1:1.7.1-2/files
%FILES%
usr/
usr/bin/
usr/bin/docker
usr/lib/
usr/lib/docker/
usr/lib/docker/dockerinit
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/docker.service
usr/lib/systemd/system/docker.socket
usr/lib/sysusers.d/
usr/lib/sysusers.d/docker.conf
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/docker
usr/share/fish/
usr/share/fish/vendor_completions.d/
usr/share/fish/vendor_completions.d/docker.fish
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/docker-attach.1.gz
usr/share/man/man1/docker-build.1.gz
usr/share/man/man1/docker-commit.1.gz
usr/share/man/man1/docker-cp.1.gz
usr/share/man/man1/docker-create.1.gz
usr/share/man/man1/docker-diff.1.gz
usr/share/man/man1/docker-events.1.gz
usr/share/man/man1/docker-exec.1.gz
usr/share/man/man1/docker-export.1.gz
usr/share/man/man1/docker-history.1.gz
usr/share/man/man1/docker-images.1.gz
usr/share/man/man1/docker-import.1.gz
usr/share/man/man1/docker-info.1.gz
usr/share/man/man1/docker-inspect.1.gz
usr/share/man/man1/docker-kill.1.gz
usr/share/man/man1/docker-load.1.gz
usr/share/man/man1/docker-login.1.gz
usr/share/man/man1/docker-logout.1.gz
usr/share/man/man1/docker-logs.1.gz
usr/share/man/man1/docker-pause.1.gz
usr/share/man/man1/docker-port.1.gz
usr/share/man/man1/docker-ps.1.gz
usr/share/man/man1/docker-pull.1.gz
usr/share/man/man1/docker-push.1.gz
usr/share/man/man1/docker-rename.1.gz
usr/share/man/man1/docker-restart.1.gz
usr/share/man/man1/docker-rm.1.gz
usr/share/man/man1/docker-rmi.1.gz
usr/share/man/man1/docker-run.1.gz
usr/share/man/man1/docker-save.1.gz
usr/share/man/man1/docker-search.1.gz
usr/share/man/man1/docker-start.1.gz
usr/share/man/man1/docker-stats.1.gz
usr/share/man/man1/docker-stop.1.gz
usr/share/man/man1/docker-tag.1.gz
usr/share/man/man1/docker-top.1.gz
usr/share/man/man1/docker-unpause.1.gz
usr/share/man/man1/docker-version.1.gz
usr/share/man/man1/docker-wait.1.gz
usr/share/man/man1/docker.1.gz
usr/share/man/man5/
usr/share/man/man5/Dockerfile.5.gz
usr/share/vim/
usr/share/vim/vimfiles/
usr/share/vim/vimfiles/ftdetect/
usr/share/vim/vimfiles/ftdetect/dockerfile.vim
usr/share/vim/vimfiles/syntax/
usr/share/vim/vimfiles/syntax/dockerfile.vim
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_docker

