# etc/systemd/grave.service
[Unit]
Description=Run sync

[Service]
User=crandel
Type=forking
ExecStart=<exec path>
