# etc/systemd/sleep.service
[Unit]
Description=Run lock screen

[Service]
User=crandel
Type=forking
Environment=DISPLAY=:0
ExecStart=/home/crandel/.sleep.fish
